# supervisord.conf
[supervisord]
nodaemon=true
loglevel=trace

[program:doorman]
environment=PYTHONPATH="/usr/src/app"
redirect_stderr=true
command=python /usr/src/app/doorman/cmd.py